# roles/certbot/files/certbot.service
[Unit]
Description=Execute certbot.sh

[Service]
Type=oneshot
ExecStart=/usr/local/bin/certbot.sh
